# build-usql-images.service
[Unit]
Description="build usql images"

[Service]
ExecStart=/home/user/src/usql/docker-images/build.sh -p
