# usql-images.service
[Unit]
Description="build usql images"

[Service]
ExecStart=/home/user/src/usql-images/build.sh -p
